# bcp-0.1/debian/postinst
#!/bin/sh
# postinst script for bcp
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        # install s3cmd from pip
        sudo pip install s3cmd
        s3cmd --configure

        if [ -f ~/.s3cfg ]; then

            echo -e "\033[42m\033[30m +++  Configure s3cmd success completed! +++ \033[0m"
            
            echo -e "\033[33m "
            cat ~/.s3cfg
            echo -e " \033[0m"

            # enable service
            sudo systemctl enable bcp.timer

            # start service
            sudo systemctl start bcp.timer

            # status service
            sudo systemctl status bcp.timer
            
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
