# bcp-0.1/debian/postinst
#!/bin/sh
# postinst script for bcp
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

S3CFGCONF=/root/.s3cfg
BCPGPGCONF=/etc/bcpgpg.conf
BCPS3CFG=/etc/bcps3cfg.conf

case "$1" in
    configure)

        #import GPG key (write conf)
        touch $BCPGPGCONF
        if [ -f $BCPGPGCONF ]; then
            echo "\033[42m\033[30m GPG keys in the system: \033[0m"
            echo "\033[33m "
            gpg --list-key
            echo " \033[0m"

            while :
            do
                read -p "Enter your GPG key (email) which can decrypt the backup archive: " email
                if [ ! -z "$email" ]; then
                    echo "$email" > $BCPGPGCONF
                    break
                fi
            done

        else
            echo "\033[31m Error: config file ${BCPGPGCONF} not exist, create a file! \033[0m"
            exit 1
        fi

        # configure s3cmd
        s3cmd --configure

        # save container name
        touch $BCPS3CFG
        if [ -f $BCPS3CFG ]; then

            while :
            do
                read -p "Enter cloud storeage container name: " container
                if [ ! -z "$container" ]; then
                    echo "$container" > $BCPS3CFG
                    break
                fi
            done
            
        else
            echo "\033[31m Error: config file ${BCPS3CFG} not exist, create a file! \033[0m"
            exit 1
        fi

        # View result
        if [ -f $S3CFGCONF ]; then

            echo "\033[42m\033[30m +++  Configure s3cmd success completed! +++ \033[0m"
            
            #echo "\033[33m "
            #cat $S3CFGCONF
            #echo " \033[0m"

            echo "\033[33m Activate the backup service by running the following commands: \033[0m"

            echo "\033[33m "
            echo "$ sudo systemctl enable bcp.timer"
            echo "$ sudo systemctl start bcp.timer"
            echo "$ sudo systemctl status bcp.timer"
            echo " \033[0m"

        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
