# bcp-0.1/bcp.service
[Unit]
Description=Backup service

[Service]
Type=oneshot
ExecStart=/usr/bin/bcp.sh

[Install]
WantedBy=multi-user.target