# bcp-0.1/debian/bcp/DEBIAN/postinst
#!/bin/sh
# postinst script for bcp
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

S3CFGCONF=/root/.s3cfg
BCPGPGCONF=/etc/bcpgpg.conf
BCPS3CFG=/etc/bcps3cfg.conf

case "$1" in
    configure)

        #import GPG key (write conf)Y
        if [ -f $BCPGPGCONF ]; then
            echo "\033[42m\033[30m GPG keys in the system: \033[0m"
            echo "\033[33m "
            gpg --list-key
            echo " \033[0m"

            read -p "Enter your GPG key (email) which can decrypt the backup archive: " email
            if [ ! -z "$email" ]; then
                echo "$email" > $BCPGPGCONF
            else
                echo "\033[31m Error: The recipient's decryption key is not specified, specify it later in the file ${$BCPGPGCONF}, otherwise it will not work! \033[0m"
                exit 1
            fi
        else
            echo "\033[31m Error: config file ${BCPGPGCONF} not exist, create a file! \033[0m"
            exit 1
        fi

        # configure s3cmd
        s3cmd --configure

        # save container name
        read -p "Enter cloud storeage container name: " container
        if [ ! -z "$container" ]; then
            echo "$container" > $BCPS3CFG
        else
            echo "\033[31m Error: The name of the container in the repository is not specified! \033[0m"
            exit 1
        fi

        if [ -f $S3CFGCONF ]; then

            echo "\033[42m\033[30m +++  Configure s3cmd success completed! +++ \033[0m"
            
            echo "\033[33m "
            cat $S3CFGCONF
            echo " \033[0m"

            # start service
            invoke-rc.d bcp.timer start >/dev/null || true

            # status service
            invoke-rc.d bcp.timer status | grep 'Active'
            
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
